# dataplane/vppagent/conf/supervisord/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug
nodaemon=true
pidfile=/run/supervisord.pid

[program:tail]
command=tail -f /dev/null
autorestart=false
redirect_stderr=true
priority=1

